# roles/mac/templates/Library/LaunchAgents/local.knakayama.homebrewupdate.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTDPLIST1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>local.knakayama.homebrewupdate</string>
        <key>ProgramArguments</key>
        <array>
          <string>/usr/local/bin/brew</string>
          <string>update</string>
        </array>
        <key>UserName</key>
        <string>{{ ansible_user_id }}</string>
        <key>StartCalendarInterval</key>
        <dict>
          <key>Hour</key>
          <integer>2</integer>
          <key>Minute</key>
          <integer>00</integer>
        </dict>
    </dict>
</plist>
