# roles/mac/templates/Library/LaunchAgents/local.knakayama.compresstmuxlog.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>local.compresstmuxlog</string>
        <key>ProgramArguments</key>
        <array>
            <string>{{ ansible_user_dir }}/bin/compress-tmux-log.rb</string>
        </array>
        <key>RunAtLoad</key>
        <true/>
        <!--
        <key>StartInterval</key>
        <integer>300</integer>
        <key>StandardOutPath</key>
        <string>/Users/knakayama/.tmuxlog/compress-tmux-log.log</string>
        <key>StandardErrorPath</key>
        <string>/Users/knakayama/.tmuxlog/compress-tmux-log.log</string>
        -->
    </dict>
</plist>
